# defaults/couchside.service
; Couchside box agent: systemd unit TEMPLATE.
; install.sh substitutes __USER__, __UID__, __EXEC__ (the resolved daemon path)
; and __CONFIG__ (the user-owned state config path) with the invoking desktop
; user's real values and installs the result to
; /etc/systemd/system/couchside.service. Do NOT hardcode /home here: the home
; can live at /var/home (Bazzite/ostree), be a systemd-homed image, or an LDAP
; path — install.sh injects the real install dir so the unit works anywhere.
[Unit]
Description=Couchside box agent
After=network-online.target
Wants=network-online.target

[Service]
User=__USER__
; input group gives the daemon access to /dev/uinput for the virtual gamepad,
; without depending on an active login seat (headless boxes have none). Paired
; with the udev rule the installer writes for /dev/uinput.
SupplementaryGroups=input
Environment=XDG_RUNTIME_DIR=/run/user/__UID__
ExecStart=/usr/bin/python3 __EXEC__ --config __CONFIG__
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
